# wall-broadcast.service
#/etc/systemd/system/
[Unit]
Description=Send message to all users

[Service]
Type=oneshot
ExecStart=/usr/local/bin/test.sh